# daily-runners.service
[Unit]
Description=Daily Runners Service
After=network.target

[Service]
Type=simple
# Change to the appropriate user for your VPS
User=nodejs
WorkingDirectory=/opt/daily-runners
ExecStart=/usr/bin/node /opt/daily-runners/dist/index.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=daily-runners
Environment=NODE_ENV=production
Environment=PORT=3000
# Make sure to set proper permissions for the .env file
EnvironmentFile=/opt/daily-runners/.env

[Install]
WantedBy=multi-user.target 